# supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=haproxy -f /etc/haproxy.cfg

[program:nodered]
command=npm start -- --userDir /data